# prog219/Week11-Final-old/final.conf
# This is an upstart script: http://upstart.ubuntu.com/index.html
description "a script to keep a node.js application in memory even after rebooting"
author      "Charle Calvert - http://www.elvenware.com/charlie"

# Start after all drives mounted
start on started mountall
stop on shutdown

# Automatically Respawn:
respawn
respawn limit 99 5

script
    export HOME="/home/ubuntu"

# The following assumes nodejs is in /usr/bin
    exec /usr/bin/nodejs /home/ubuntu/bin/final/bin/www >> /var/log/final.log 2>&1

end script

post-start script
   # Optionally put a script here that will notifiy you node has (re)started
   # /root/bin/hoptoad.sh "node.js has started!"
   echo "node.sj has started running BridgeReader"
end script
